# power-logger.service.template
[Unit]
Description=Power stats CSV logger
After=default.target

[Service]
Type=simple
ExecStart=@CURDIR@/venv/bin/python @CURDIR@/record.py
Restart=on-failure

[Install]
WantedBy=default.target
